# systemd/mb-media-to-talk.service
[Unit]
Description=Matterbridge media relay to Nextcloud Talk
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/mb-media-to-talk
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
